# packages/filebeat/extra/filebeat.service
[Unit]
Description=Filebeat log forwarder

[Service]
Restart=always
StartLimitInterval=0
RestartSec=5
KillMode=control-group
Delegate=true
LimitNOFILE=infinity
ExecStartPre=${PKG_PATH}/bin/filebeat_kafka_collector_hosts.sh
ExecStart=${PKG_PATH}/bin/filebeat/filebeat -c ${PKG_PATH}/etc/filebeat.yml
