# packages/grafana/extra/grafana.service
[Unit]
Description=Grafana 
# if external kafka collector is provided, disable collector stack (kafka,ingestion,telegraf-collector,influxdb,elasticsearch,grafana)
ConditionFileNotEmpty=!/opt/mesosphere/etc/kafka_collector_hosts
[Service]
Restart=always
StartLimitInterval=0
RestartSec=5
KillMode=control-group
Delegate=true
LimitNOFILE=infinity
ExecStart=${PKG_PATH}/bin/grafana-server --homepath ${PKG_PATH}/
